# debian/pi-relay-control/lib/systemd/system/pi-relay-control.service
[Unit]
Description=Pi Relay Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/relay_control
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
